# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /data/postgres
user=postgres
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

[program:seaweedfs]
command=/usr/local/bin/weed server -filer -s3 -s3.port=8333 -s3.config=/etc/seaweedfs/s3.json -volume.port=8180 -dir=/data/seaweedfs -master.dir=/data/seaweedfs/master -volume.dir.idx=/data/seaweedfs -ip=localhost
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20

[program:trino]
command=/opt/trino/bin/launcher run
user=trino
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=30
